# push-to-talk/ptt.service
[Unit]
Description=Push to talk service

[Service]
ExecStart=/usr/local/bin/ptt.py

[Install]
WantedBy=multi-user.target
